# data/backups.service
[Service]
Type=oneshot
ExecStart=/srv/backups/run.sh

